# supervisord.conf
[supervisord]
nodaemon=true

[program:openduty]
command=bash -c "/opt/openduty/start_openduty.sh"
